# reddit_posts/crontab
30 20 * * * python3 reddit_posts/scrape.py >> /log
